# cfusion/runtime/bin/wsconfig
#!/bin/sh

# chkconfig: 345 90 14
# description:ColdFusion Server connector configuration

CFUSION_HOME=$(cd "$(dirname "$0")"; pwd)/../..
. $CFUSION_HOME/bin/findjava.sh
JVM_ARGS="-server -Xms256m -Xmx512m"

PATH=/usr/xpg4/bin:/bin:/sbin:/usr/bin:/usr/sbin:$PATH
CONNECTOR=""
IS_RUNNING="false"
IS_MAC="false"
VERSION="ColdFusion 2021"

CLASSPATH=$CFUSION_HOME/runtime/lib/wsconfig.jar

CLASSPATH=$CLASSPATH:$JAVA_CLASSPATH
JAVA_LIBRARY_PATH=$JAVA_LIBRARYPATH
JAVA_EXECUTABLE=$_JAVACMD

ID=`id -u`

if [ "`uname`" = "Darwin" ]; then
  # nothing, on OSX we don't care about the user name.  We assume that the user running it has permission to do everything
  IS_MAC="true" 
elif [ ! $ID -eq 0 ]; then
	echo "You must be root to invoke wsconfig."
	exit 1
fi

wsconfig() {

	eval $CFSTART

}

case `uname` in

	SunOS)
		OS=Solaris
		LD_LIBRARY_PATH="$CFUSION_HOME/lib:$CFUSION_HOME/lib/_solaris/bin:$JAVA_LIBRARY_PATH"
                ARGS=$*
		CFSTART='LD_LIBRARY_PATH=$LD_LIBRARY_PATH; export LD_LIBRARY_PATH; cd $CFUSION_HOME/bin; $JAVA_EXECUTABLE -classpath $CLASSPATH $JVM_ARGS com.adobe.coldfusion.connector.connectorinstaller.ConnectorInstaller $ARGS'
	;;

	Darwin)
		OS=Darwin
		LD_LIBRARY_PATH="$CFUSION_HOME/lib:$JAVA_LIBRARY_PATH"
		ARGS=$*
		CFSTART='LD_LIBRARY_PATH=$LD_LIBRARY_PATH; export LD_LIBRARY_PATH;cd $CFUSION_HOME/bin; $JAVA_EXECUTABLE -classpath $CLASSPATH $JVM_ARGS com.adobe.coldfusion.connector.connectorinstaller.ConnectorInstaller $ARGS'
	;;

	Linux)
		OS=Linux
		LD_LIBRARY_PATH="$CFUSION_HOME/lib:$CFUSION_HOME/lib/_ilnx21/bin:$JAVA_LIBRARY_PATH"
		ARGS=$*
		CFSTART='export LD_LIBRARY_PATH=$LD_LIBRARY_PATH; cd $CFUSION_HOME/bin; $JAVA_EXECUTABLE -classpath $CLASSPATH $JVM_ARGS com.adobe.coldfusion.connector.connectorinstaller.ConnectorInstaller $ARGS'

	;;

	*)
		echo "Your OS: `uname` is unsupported"
		echo "exiting"
		exit 1
	;;

esac

ARG=$1

[ -z "$ARG" ] && ARG=usage

wsconfig

exit 0

